# resources/mako-master.conf
description     "Start Mako Master"
author          "@travellperkins"

start on runlevel [2345]
stop on runlevel [2345]

respawn
respawn limit 5 60

setuid mako
setgid mako

env NODE_BIN=/usr/local/node

script
        chdir /apps/mako
        . ./.env
        export MINION_PASS
        export MASTER_HOST
        export MASTER_PASS
        PATH=$NODE_BIN:$PATH
        exec /usr/bin/env field-marshal 2>&1 | /usr/bin/env logger -t mako-host
end script
